# crontab
@reboot echo $(date) >> rebootDateTimes.txt && /home/pi/start.cron &
